# crontab
30 8,13,16,20 * * 1-5 /home/pi/dingding.sh > /dev/null 2>&1
*/30 9-20 * * 1-5 /home/pi/kua.sh > /dev/null 2>&1

